# salt/monitoring/vmagent/templates/vmagent.service.jinja
[Unit]
Description=VictoriaMetrics Agent
After=network.target

[Service]
Type=simple
User=vmagent
Group=vmagent
ExecStart=/usr/sbin/vmagent \
  -promscrape.config=/etc/vmagent/scrape.yml \
  -remoteWrite.url={{ vmagent.get('remote_write_url', 'http://localhost:8428/api/v1/write') }} \
  -httpListenAddr=:8429
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
